# aii/deploy/systemd/aii-pg-backup.timer
[Unit]
Description=Daily AII Postgres backup at 03:15

[Timer]
OnCalendar=*-*-* 03:15:00
Persistent=true
AccuracySec=5min

[Install]
WantedBy=timers.target
