# aii/deploy/systemd/aii-context-archiver.timer
[Unit]
Description=AII context-archiver 每小时归档一轮(alldata 无文件时秒退)

[Timer]
OnCalendar=hourly
RandomizedDelaySec=120
Persistent=true

[Install]
WantedBy=timers.target
