# aii/deploy/systemd/aii-watchdog.timer
[Unit]
Description=Run AII watchdog every 5 minutes

[Timer]
OnBootSec=90
OnUnitActiveSec=300
AccuracySec=15

[Install]
WantedBy=timers.target
